# ops/debian/etc/init/recall-server.conf
# Ubuntu upstart file

description "Recall bookmark web service - server"
author "Cal Paterson <cal@calpaterson.com>"

start on runlevel [2345]
stop on runlevel [06]

exec /usr/bin/recall-server